# crontab
*/1 * * * * root /usr/local/bin/python /app/main.py >> /var/log/cron/cron.log 2>&1
